# templates/etc/profile.d/linuxmuster-client.sh
# linuxmuster-client profile
if [ -e /usr/share/linuxmuster-client/profile ]; then
    . /usr/share/linuxmuster-client/profile
fi

umask 022
